# discoverd/systemd/nvme-discoverd.service.in
# SPDX-License-Identifier: GPL-2.0-or-later
[Unit]
Description=NVMe-oF Discovery Daemon
Documentation=man:nvme-discoverd(8)
Before=remote-fs.target

[Service]
Type=notify-reload
ExecStart=@SBINDIR@/nvme-discoverd
Restart=on-failure
RuntimeDirectory=nvme/discoverd nvme/discoverd/units nvme/discoverd/controllers
RuntimeDirectoryPreserve=yes

[Install]
WantedBy=multi-user.target
