# nvmf-autoconnect/systemd/nbft-boot-connect.service.in
# Establish NVMe/TCP connections as defined in the ACPI NBFT table

[Unit]
# Prevent deadlock by avoiding systemd adding dependencies
# on unreachable targets.
DefaultDependencies=no
# Specify 'rd.nvmf.nonbft' to disable NBFT boot
ConditionKernelCommandLine=!rd.nvmf.nonbft
ConditionPathExistsGlob=/sys/firmware/acpi/tables/NBFT*
Wants=modprobe@nvme_fabrics.service
After=modprobe@nvme_fabrics.service
Wants=network-online.target
After=network-online.target
Before=remote-fs-pre.target

[Service]
RestrictAddressFamilies=AF_INET AF_INET6
Type=oneshot
# Connection attempt recovery
Restart=on-failure
RestartSec=10
StartLimitIntervalSec=0
StartLimitBurst=0
ExecStart=/usr/sbin/nvme connect-all --nbft

[Install]
WantedBy=sysinit.target
